# rootfs/usr/lib/systemd/system/tornado.service
[Unit]
    Description=Start all snort instances of tornado
[Service]
    Type=oneshot
    ExecStart=/bin/true
    RemainAfterExit=yes

[Install]
    WantedBy=multi-user.target
