# rootfs/usr/lib/systemd/system/tornado.service
[Unit]
    Description=Tornado supervisor
[Service]
    Type=oneshot
    ExecStart=/srv/http/pychat/backend/.venv/bin/python /srv/http/pychat/backend/manage.py flush_online
    RemainAfterExit=yes
    User=http

[Install]
    WantedBy=multi-user.target
