# rootfs/usr/lib/systemd/system/tornado.service
[Unit]
    Description=Tornado supervisor
[Service]
    Type=oneshot
    ExecStart=/home/alarm/djangochat/manage.py flush_online
    RemainAfterExit=yes
		User=http

[Install]
    WantedBy=multi-user.target
