# docker-all/tornado
#!/sbin/openrc-run

command="python3"
command_args="/srv/http/manage.py start_tornado --host localhost"
command_background="yes"
command_user="http"
pidfile="/tmp/tornado.pid"

depend() {
    want mariadb redis
}

start_pre() {
   su - http -c 'PYCHAT_CONFIG=docker_all python3 /srv/http/manage.py compile_js'
   su - http -c 'PYCHAT_CONFIG=docker_all python3 /srv/http/manage.py init_db'
   su - http -c 'PYCHAT_CONFIG=docker_all python3 /srv/http/manage.py sync_db'
}