# docker-all/tornado
#!/sbin/openrc-run

command="python3"
command_args="/srv/http/manage.py init_db_start_tornado"
command_background="yes"
pidfile="/tmp/tornado.pid"

depend() {
    want mariadb redis
}